# vendor/github.com/nlamirault/abraracourcix/contrib/systemd/abraracourcix.service
[Unit]
Description=Abraracourcix

[Service]
ExecStart=/usr/bin/abraracourcix /etc/abraracourcix.toml
Restart=on-failure